# deploy/mylibpro.service
[Unit]
Description=MyLibPro (Next.js academic library)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Source checkout on the server (see deploy/README.md).
WorkingDirectory=/opt/mylibpro/app
# Runtime config: API keys + DATA_ROOT/DB_PATH (host paths, not container paths).
EnvironmentFile=/opt/mylibpro/.env
Environment=NODE_ENV=production
# Bind to loopback only — host Caddy terminates TLS and proxies to it.
# (`next start` reads -H/-p flags, not the HOSTNAME/PORT env vars.)
ExecStart=/opt/mylibpro/app/node_modules/.bin/next start -H 127.0.0.1 -p 3000
Restart=always
RestartSec=5
# Light hardening for a single-app box.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
